# systemd/chatterbox-turbo-fastapi.service
[Unit]
Description=Chatterbox Turbo FastAPI API
After=network-online.target
Wants=network-online.target
PartOf=chatterbox-turbo.target

[Service]
Type=simple
User=op
Group=op
WorkingDirectory=/home/op/chatterbox-turbo-fastapi
EnvironmentFile=-/etc/chatterbox-turbo-fastapi.env
ExecStart=/bin/bash /home/op/chatterbox-turbo-fastapi/run_api_service.sh
Restart=always
RestartSec=5
TimeoutStopSec=30
KillSignal=SIGTERM
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
